# cronjob-clamav.service.in
[Unit]
Description=clamav cronjob

[Service]
ExecStart=@sbindir@/clamav-cron-updatedb
