# cronjob-clamav.timer
[Unit]
Description=clamav cronjob timer

[Timer]
OnCalendar=*:05

[Install]
WantedBy=cronjobs.target
